# ci/09-ci-05-teamcity/scripts/teamcity.service
[Unit]
Description=TeamCity Server and Agent (Docker Compose)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/start-teamcity.sh
ExecStop=/usr/bin/docker compose down

[Install]
WantedBy=multi-user.target